# roles/acpi/files/etc/udev/rules.d/99-presentation.rules
# Managed by Ansible

SUBSYSTEM=="drm", ACTION=="change", RUN+="/etc/acpi/presentation.sh"
